# .config/systemd/user/paperChanger.timer
[Unit]
Description=Run Paperchanger every 2 min

[Timer]
OnBootSec=60sec
OnUnitActiveSec=3min

[Install]
WantedBy=timers.target
