# .config/systemd/user/paperChanger.service
[Unit]
Description=Randomly selects a wallpaper

[Service]
Type=simple
ExecStart=%h/bin/paperChanger
Restart=no
